# src/applications/mpt-daemons/sender/config/mpt-sender-daemon.service.in
[Unit]
Description=MPT Sender Daemon
After=syslog.target network.target

[Service]
Type=forking
EnvironmentFile=-/@CMAKE_INSTALL_SYSCONFIG_PATH@/mpt-sender-daemon
ExecStart=@CMAKE_INSTALL_PREFIX@/bin/mpt-sender-daemon -m $MAESTRO_BROKER_URL -L $MPT_LOG_DIR -l $MPT_LOG_LEVEL -n $NODE_NAME
PrivateTmp=false

[Install]
WantedBy=multi-user.target